# systemd/proxmox-ha-hwmon.service
[Unit]
Description=Publish Proxmox hardware sensors to Home Assistant via MQTT
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxmox-ha-hwmon-mqtt.sh
