# systemd/proxmox-ha-smart.service
[Unit]
Description=Publish Proxmox disk SMART/NVMe data to Home Assistant via MQTT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/proxmox-ha-smart-daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
